# caches/skills/wsl-bypass-vpn/scripts/wsl-bypass-vpn.service
[Unit]
Description=Bypass VPN for WSL mirrored networking
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wsl-bypass-vpn.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
