# caches/skills/wsl-bypass-vpn/scripts/wsl-bypass-vpn.timer
[Unit]
Description=Re-apply WSL VPN bypass routes after network changes

[Timer]
OnBootSec=20
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
